# com.dcselfie.arca-tunnel.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.dcselfie.arca-tunnel</string>

    <!--
      맥이 OCI를 향해 리버스 SSH 연결을 열고, "-R 1080"(목적지 미지정)으로
      OCI:1080에 도달하는 연결을 맥 자신의 SOCKS 4/5 프록시로 서비스한다.
      (ssh(1) man: "if no explicit destination was specified, ssh will act
      as a SOCKS 4/5 proxy and forward connections to the destinations
      requested by the remote SOCKS client.")

      즉 실제 목적지(arca.live) 연결은 맥이 직접 연다 → 맥의 공인 IP로 나감.

      이전 버전은 "-D 1080 -R 1080:localhost:1080"을 같이 썼는데, 이러면
      -R이 맥의 localhost:1080(=-D가 만든, 출구가 OCI인 SOCKS)으로 되돌아가
      결국 다시 OCI를 거쳐 나가는 루프가 되어 맥 IP가 전혀 쓰이지 않았다.
      -D를 빼고 -R 단독(목적지 없음)으로 리버스 다이나믹 포워딩만 쓴다.

      OCI sshd 쪽은 기본값(loopback 바인딩)이라 GatewayPorts 등 별도 설정 불필요.
    -->
    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/autossh</string>
        <string>-M</string>
        <string>0</string>
        <string>-N</string>
        <string>-o</string>
        <string>ServerAliveInterval=30</string>
        <string>-o</string>
        <string>ServerAliveCountMax=3</string>
        <string>-o</string>
        <string>ExitOnForwardFailure=yes</string>
        <string>-R</string>
        <string>1080</string>
        <string>ubuntu@158.180.94.136</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/tmp/arca-tunnel.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/arca-tunnel.log</string>
</dict>
</plist>
