# dcselfie-web.service
[Unit]
Description=dcselfie web gallery
After=network-online.target
Before=dcselfie-launcher.service
Wants=network-online.target

[Service]
User=ubuntu
WorkingDirectory=/home/ubuntu/dcinsideImageCrawler
Environment=WEB_HOST=127.0.0.1
Environment=WEB_PORT=8000
ExecStart=/home/ubuntu/dcinsideImageCrawler/venv/bin/python run_web_server.py
Restart=on-failure
RestartSec=5
LimitCORE=0
MemoryMax=512M
MemorySwapMax=0
UMask=0077

[Install]
WantedBy=multi-user.target
