# dcselfie-web.service
[Unit]
Description=dcselfie web gallery
After=network-online.target
Before=dcselfie-launcher.service
Wants=network-online.target

[Service]
User=ubuntu
WorkingDirectory=/home/ubuntu/dcinsideImageCrawler
# systemd가 .env를 직접 읽도록 해 WEB_INGEST_TOKEN이 웹 프로세스에 항상 전달된다.
# '-'는 초기 설치처럼 아직 .env가 없는 경우에도 서비스 등록 자체는 허용한다.
EnvironmentFile=-/home/ubuntu/dcinsideImageCrawler/.env
Environment=WEB_HOST=127.0.0.1
Environment=WEB_PORT=8000
ExecStart=/home/ubuntu/dcinsideImageCrawler/venv/bin/python run_web_server.py
Restart=on-failure
RestartSec=5
LimitCORE=0
# 스토어 상한(기본 256MB) + PIL 디코드 스파이크(대형 이미지 사본 포함 ~72MB × 동시 2슬롯)
# + 인터프리터(~100MB)를 감안한 여유. 512M는 만충 상태에서 cgroup OOM(SIGKILL) 여지가 있었다.
MemoryMax=768M
MemorySwapMax=0
UMask=0077

[Install]
WantedBy=multi-user.target
