# dcselfie-launcher.service
[Unit]
Description=dcselfie image crawler launcher (DCInside + Arcalive)
Requires=dcselfie-web.service
After=network.target network-online.target dcselfie-web.service
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/dcinsideImageCrawler
Environment=WEB_GALLERY=1
Environment=PYTHONUNBUFFERED=1
# ARCA_SOCKS_PROXY는 여기(git 추적 파일)에 두지 않는다 — 서버의 .env에만 둔다
# (launcher.py가 이미 load_dotenv()로 .env를 읽으므로 systemd Environment= 불필요).
# 예: echo 'ARCA_SOCKS_PROXY=socks5://localhost:1080' >> .env
ExecStart=/home/ubuntu/dcinsideImageCrawler/venv/bin/python3 launcher.py
Restart=always
RestartSec=30
LimitCORE=0
UMask=0077
StandardOutput=append:/home/ubuntu/dcinsideImageCrawler/logs/launcher.log
StandardError=append:/home/ubuntu/dcinsideImageCrawler/logs/launcher.log

[Install]
WantedBy=multi-user.target
